# deploy/supervisor.conf
[program:agent]
command=python backend/app/agent.py
directory=/Users/mac/Documents/ai/learnclaudecode/learnclaudecode
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/tmp/agent.log
stderr_logfile=/tmp/agent_error.log
user=mac
environment=PYTHONUNBUFFERED=1

; 停止信号
stopsignal=TERM
stopwaitsecs=10

; 重启条件
autorestart=unexpected
exitcodes=0

; 健康检查（需要 supervisor 3.0+）
; 如果进程退出码不是 0，自动重启
